# wifi_watch.conf
[program:wifi_watch]
command=/root/src/watch.sh
user=wifi_watch
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/root/src/logs/wifi_watch.log
